# iptracking-daemon.service.in
[Unit]
Description=PAM SSH event logging daemon
After=network.target
Before=sshd.service

[Service]
Type=simple
EnvironmentFile=-/etc/sysconfig/iptracking-daemon
ExecStart=@CMAKE_INSTALL_FULL_SBINDIR@/iptracking-daemon $OPTIONS

[Install]
WantedBy=multi-user.target