# pam-daemon/iptracking-pamd.service.in
[Unit]
Description=PAM SSH event logging daemon
After=network.target
Before=sshd.service

[Service]
Type=simple
EnvironmentFile=-/etc/sysconfig/iptracking-pamd
ExecStart=@CMAKE_INSTALL_FULL_SBINDIR@/iptracking-pamd $OPTIONS

[Install]
WantedBy=multi-user.target